# #dev/supervisor/conf.d/vbox/server-auth-application-dev-daemon-redis-queue.ini
[program:server-auth-application-dev-daemon-redis-queue]
command=php /var/www/home-www/yusam/github/yusam-hub/project-0001/server-auth-application/console daemon:redis-queue
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor/server-auth-application-dev-daemon-redis-queue.log
user=apache
autostart=false
autorestart=false